# mnemos.service
[Unit]
Description=MNEMOS API Server - Memory System with Integrated Compression
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
# Change to the user that owns /opt/mnemos
User=mnemos
WorkingDirectory=/opt/mnemos
Environment="PATH=/opt/mnemos/venv/bin"
Environment="MNEMOS_BIND=0.0.0.0"
Environment="MNEMOS_PORT=5002"
EnvironmentFile=/opt/mnemos/.env
# Canonical port is 5002; this overrides [api].port in config.toml if different
# MNEMOS v3.5 is single-worker by design; do not configure multiple Uvicorn workers.
# Keep the venv path explicit so systemd uses the installed console script.
ExecStart=/opt/mnemos/venv/bin/mnemos serve

# Auto-restart on failure
Restart=on-failure
RestartSec=10s
StartLimitBurst=5

# Resource limits
MemoryMax=2G
CPUQuota=80%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mnemos

# Security
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
